# templates/example-system.service
# /etc/systemd/system/example-system.service

[Unit]
Description=Example System Service (e.g., Web Server)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/nginx -g "daemon off;"
Restart=on-failure
RestartSec=5s
User=http
Group=http
Environment="VAR=value"
WorkingDirectory=/etc/nginx

[Install]
WantedBy=multi-user.target